# install/templates/etc/crontabs/root
*/5 * * * * /root/admin/isAlive.sh
0 3 * * * reboot 
